# backend/pipeline/crontab
0 4 * * * cd /app/backend && python -m pipeline.run
