# backend/pipeline/crontab
0 4 * * * cd /app/backend && /app/.venv/bin/python -m pipeline.run
